# crontab
*/2 * * * * /bin/convert >> /var/log/cron.log 2>&1

# Mandatory blank line
